# service/hive-indexer/run
#!/bin/bash

POPULATE_CMD="$(which hive)"

exec python3.5 -u "${POPULATE_CMD}" indexer run
