# service/hive-indexer/run
#!/bin/bash

POPULATE_CMD="$(which hive)"

exec python3.5 "${POPULATE_CMD}" indexer from-steemd \
#  --database_url "${DATABASE_URL:-sqlite:///}" \