# service/hive-indexer/run
#!/bin/bash

POPULATE_CMD="$(which hive)"

exec "${POPULATE_CMD}" indexer run
