# service/hive-web/run
#!/bin/bash

POPULATE_CMD="$(which hive)"

python3 /app/hive/run.py server --database-url="${DATABASE_URL}" --http-server-port="${HTTP_SERVER_PORT}"
